# server/systemd/alpaca-zero-dte-exit-review.timer
[Unit]
Description=Run the paper 0DTE Level 2 exit review every minute during market hours

[Timer]
OnCalendar=Mon..Fri *-*-* 09..15:*:55
Persistent=false
RandomizedDelaySec=5
AccuracySec=1s
Unit=alpaca-zero-dte-exit-review.service

[Install]
WantedBy=timers.target
