# server/systemd/alpaca-zero-dte-engine.timer
[Unit]
Description=Run the paper 0DTE Level 2 engine every minute during the entry window

[Timer]
OnCalendar=Mon..Fri *-*-* 09:35/1:45
OnCalendar=Mon..Fri *-*-* 10..14:*:45
OnCalendar=Mon..Fri *-*-* 15:00..15:45
Persistent=false
RandomizedDelaySec=5
AccuracySec=1s
Unit=alpaca-zero-dte-engine.service

[Install]
WantedBy=timers.target
