# server/systemd/alpaca-zero-dte-reconcile.timer
[Unit]
Description=Reconcile paper 0DTE Level 2 state every five minutes during market hours

[Timer]
OnCalendar=Mon..Fri *-*-* 09..15:1/5:30
Persistent=false
RandomizedDelaySec=5
AccuracySec=1s
Unit=alpaca-zero-dte-reconcile.service

[Install]
WantedBy=timers.target
