# server/systemd/alpaca-universe-lifecycle.service
[Unit]
Description=Alpaca Autonomous Universe Lifecycle Service
After=network-online.target alpaca-zero-dte-eod.service paper-ops-late-day.service
Wants=network-online.target

[Service]
Type=oneshot
User=alpaca
Group=alpaca
WorkingDirectory=/home/alpaca/Alpaca-Trading
EnvironmentFile=/opt/alpaca-investing/secrets/alpaca.env
Environment=PATH=/usr/bin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=NODE_ENV=production
Environment=AUTOMATED_PAPER_EXECUTION_ENABLED=false
ExecStart=/usr/bin/npm run universe:lifecycle -- --format=json
TimeoutStartSec=900
SyslogIdentifier=alpaca-universe-lifecycle
StandardOutput=journal
StandardError=journal
