# server/systemd/alpaca-zero-dte-eod.timer
[Unit]
Description=Create the paper 0DTE Level 2 end of day summary after the force exit window

[Timer]
OnCalendar=Mon..Fri *-*-* 16:05:00
Persistent=false
RandomizedDelaySec=15
AccuracySec=1s
Unit=alpaca-zero-dte-eod.service

[Install]
WantedBy=timers.target
