# server/systemd/paper-ops-late-day.timer
[Unit]
Description=Run Alpaca Paper Ops Late-Day Review at 3:25 PM ET on weekdays

[Timer]
OnCalendar=Mon..Fri *-*-* 15:25:00
Persistent=false
Unit=paper-ops-late-day.service

[Install]
WantedBy=timers.target
