# server/systemd/paper-ops-midday.timer
[Unit]
Description=Run Alpaca Paper Ops Midday Portfolio Review at 12:10 PM ET on weekdays

[Timer]
OnCalendar=Mon..Fri *-*-* 12:10:00
Persistent=false
Unit=paper-ops-midday.service

[Install]
WantedBy=timers.target
